# dev-db/cassandra/files/init_d_cassandra
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

PIDFILE=/var/lib/cassandra/cassandra.pid

depend() {
	need net
}

start() {
	source /usr/share/cassandra/cassandra.in.sh
	ebegin "Starting cassandra"
	start-stop-daemon --start \
		--chuid cassandra:cassandra \
		--pidfile ${PIDFILE} --make-pidfile \
		--exec /usr/bin/java -- \
		${JVM_OPTS} -Dstorage-config=${CASSANDRA_CONF} \
		-Dcassandra-pidfile=${PIDFILE} \
		-Dlog4j.configuration=log4j-server.properties \
		-cp "${CLASSPATH}:/etc/cassandra" \
		org.apache.cassandra.thrift.CassandraDaemon <&- > /dev/null &
	eend $?
}

stop() {
	ebegin "Stopping cassandra-daemon"
	start-stop-daemon --stop --quiet \
		--pidfile ${PIDFILE}
	eend $?
}
